# service/deskworks.service.example
# Always-on Deskworks on Linux (systemd user service).
#
#   1. Edit the REPLACE_ME paths.
#   2. Install:
#        mkdir -p ~/.config/systemd/user
#        cp service/deskworks.service.example ~/.config/systemd/user/deskworks.service
#        systemctl --user daemon-reload
#        systemctl --user enable --now deskworks
#   3. Open http://127.0.0.1:5007 . Logs:  journalctl --user -u deskworks -f
#   4. Stop:  systemctl --user disable --now deskworks

[Unit]
Description=Deskworks — private local AI over your documents
After=network.target

[Service]
# REPLACE_ME: folder containing your deskworks.toml
WorkingDirectory=REPLACE_ME
# REPLACE_ME: python in your venv
ExecStart=REPLACE_ME/.venv/bin/python -m deskworks.cli web
Environment=DESKWORKS_EMBED_DEVICE=cpu
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
